# sway/.config/systemd/user/sway-bg.service
[Unit]
Description=Select a new background image

[Service]
ExecStart=%h/.config/sway/select-bg.py

[Install]
WantedBy=multi-user.target
